# apps/courier/deploy/roles/courier/templates/chelix-courier.service.j2
[Unit]
Description=Chelix Courier – APNS push relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=chelix-courier
Group=chelix-courier
EnvironmentFile=/etc/chelix-courier/env
ExecStart=/usr/local/bin/chelix-courier \
    --bind 127.0.0.1 \
    --port {{ courier_port }} \
    --key-path /etc/chelix-courier/AuthKey.p8 \
    --key-id ${COURIER_KEY_ID} \
    --team-id ${COURIER_TEAM_ID} \
    --bundle-id ${COURIER_BUNDLE_ID} \
    --auth-token ${COURIER_AUTH_TOKEN}
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/etc/chelix-courier

[Install]
WantedBy=multi-user.target
